# supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/bin/bash /start-spark.sh

[program:flask]
command=python3 /home/app.py